# config/listen.service
[Unit]
Description=Mini-Agent Listen Server
After=network.target graphical.target

[Service]
Type=simple
WorkingDirectory=/home/jd/mac-mini-agent
ExecStart=/home/jd/mac-mini-agent/config/start-listen.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
